# crontab/get_volume_data.sh
#! /bin/bash

echo "Extraction des données J-1 sur Binance, pour les données suivantes :"
echo "Open price, High price, Low price, Close price, Volume, Number of Trades, Taker assets buy volume"

python3 /home/sanou/BTC/1_get_data_scripts/get_binance_daily_data.py